# scripts/gabi-sync@.service
[Unit]
Description=GABI DOU Daily Sync — catalog refresh + download + ingest + BM25
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=%i
WorkingDirectory=/home/%i/dev/gabi-kimi
ExecStart=/home/%i/dev/gabi-kimi/scripts/daily_sync.sh
TimeoutStartSec=7200
StandardOutput=journal
StandardError=journal

# Prevent overlapping runs
ExecStartPre=/bin/bash -c 'pgrep -f "ingest.sync_pipeline" && echo "Sync already running" && exit 1 || true'

[Install]
WantedBy=multi-user.target
